# supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini #APP_ROOT#/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

[supervisord]